# docker/rootfs/etc/s6-overlay/s6-rc.d/ssh-init/up
#!/command/execlineb -P
sh -c /etc/s6-overlay/scripts/configure-ssh.sh