# debian/prerm
#!/bin/sh
set -e
# Pre-removal: stop and disable service
if command -v systemctl >/dev/null 2>&1; then
  systemctl stop modbus2mqtt-server.service || true
  systemctl disable modbus2mqtt-server.service || true
  systemctl daemon-reload || true
fi
exit 0
