# docker/rootfs/etc/s6-overlay/s6-rc.d/modbus2mqtt/run
#!/bin/sh
set +e
cd /usr/app
npx modbus2mqtt --yaml /data --ssl /ssl
