# docker/rootfs.debug/etc/s6-overlay/s6-rc.d/modbus2mqtt/run
#!/command/with-contenv bashio
set +e
APPDIR=/usr/app
cd $APPDIR
DEBUG_OPT="--node-options=--inspect-brk=0.0.0.0 "
DEBUG="m2mgithub config.addon"
TERMINATE="false"
export HASSIO_TOKEN
while [ "$TERMINATE" = "false" ]
do
   bashio::log.info "Starting  starting modbus2mqtt \"$HASSIO_TOKEN\""
   exec npx $DEBUG_OPT modbus2mqtt --yaml /data --ssl /ssl --term 0
done